# deployment/mininter-scheduler.service
[Unit]
Description=MININTER GPS Scheduler
After=network.target mysql.service
Requires=network.target

[Service]
Type=simple
User=www-data
Group=www-data
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Rutas del proyecto (CAMBIAR POR TU RUTA DE PRODUCCIÓN)
WorkingDirectory=/var/www/mininter
ExecStart=/usr/bin/php /var/www/mininter/artisan schedule:work

# Variables de entorno
Environment=LARAVEL_ENV=production
Environment=PHP_INI_SCAN_DIR=/etc/php/8.2/cli/conf.d

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mininter-scheduler

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/mininter/storage /var/www/mininter/bootstrap/cache

# Proceso hijo
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target 